# k8s-node/container/etc/systemd/system/kubeconfig.path
[Unit]
Description=Watch for kubeconfig
[Path]
PathExists=/mnt/share/data/kubeconfig.node
[Install]
WantedBy=multi-user.target
