# k8s-master/container/etc/.overlay/master/systemd/system/kubelet.service.d/10-master-flags.conf
[Service]
Environment="KUBELET_ARGS_EXTRA=\
        --node-labels=node-role.kubernetes.io/master,master=true \
    "
