# k8s-node/container/etc/systemd/system/docker.service.d/10-overlay.conf
[Service]
ExecStart=
ExecStart=/usr/bin/dockerd -H fd:// --storage-driver=overlay2 --storage-opt overlay2.override_kernel_check=true
