# data/overlay/systemd/system/kubelet.service.d/10-fail-on-swap.conf
[Service]
Environment="KUBELET_ARGS=\
        --fail-swap-on=false \
    "
